# web/client/systemd/farmpulse-sidecar.service
[Unit]
Description=FarmPulse sidecar (send stats to FarmPulse; does not touch Hive cloud)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/farmpulse/bin/sidecar.sh
User=root
# Иначе после выхода sidecar.sh systemd может завершить дочерний reboot/sreboot в той же cgroup
KillMode=none
