# web/client/systemd/farmpulse-sidecar.service
[Unit]
Description=FarmPulse sidecar (send stats to FarmPulse; does not touch Hive cloud)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Иначе при выходе sidecar systemd убивает весь cgroup — фоновый /hive/sbin/sreboot не доживает до ребута.
KillMode=process
ExecStart=/opt/farmpulse/bin/sidecar.sh
User=root
SyslogIdentifier=farmpulse-sidecar
StandardOutput=journal
StandardError=journal
